# focus-bot.service
[Unit]
Description=Focus Bot - Telegram to Obsidian note capture
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%WORKING_DIR%
ExecStart=%BUN_PATH% run start
Restart=always
RestartSec=5
EnvironmentFile=%WORKING_DIR%/.env

[Install]
WantedBy=multi-user.target
